# systemd/syncra-ws-gateway.service
[Unit]
Description=Syncra: Superfast Local Model Chat UI
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/syncra
EnvironmentFile=/etc/syncra-ws-gateway.env
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 9000
Restart=on-failure
RestartSec=2
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
